# target/linux/allwinner6.3/base-files/etc/uci-defaults/51-unmute-sound
amixer -c 0 set Headphone 100% unmute

exit 1